# recipes-calaos/calaos-web/calaos-web/calaos-update-config.service

[Unit]
Description=Update calaos configuration
DefaultDependencies=no

[Service]
ExecStart=/usr/bin/php /www/pages/config_update.php

